# scfw.service
[Unit]
Description=iptables firewall

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'sh /usr/lib/scfw.sh'

[Install]
WantedBy=multi-user.target
