# scfw.service
[Unit]
Description=iptables firewall

[Service]
Type=oneshot
ExecStart=/usr/bin/restartscfw

[Install]
WantedBy=multi-user.target
